# overlay/etc/systemd/system/ubuntu-lte.service
#start
[Unit]
Description=Setup ubuntu platform environment of 4g

[Service]
Type=simple
ExecStart=/etc/init.d/ubuntu-lte.sh

#end

