# overlay/etc/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="udc",ACTION=="change",DRIVER=="configfs-gadget",RUN+="/usr/local/bin/usbdevice %E{DEVPATH}"
